# net-misc/sia/files/initd/siad
#!/sbin/openrc-run
depend() {
  need localmount net
  after bootmisc
  use logger
}
extra_commands="fstart fstop"

start() {
	ebegin "Starting siad daemon"
	# Sourcing parameters.txt
	/usr/local/bin/paramcopy
	# We check if ${SIAUSR} exists. If not, create it.
	if id "${SIAUSR}" >/dev/null 2>&1; then
                :
	else
	       /usr/sbin/useradd -m -G users,wheel,audio -s /bin/bash ${SIAUSR}
	fi
	if [ ! -d "${SCDIR}" ]; then
	  mkdir -p ${SCDIR}
	fi
	chown ${SIAUSR}:${SIAUSR} -R ${SCDIR}
	mpoint="$(dirname "${SCDIR}")"
	swapon -p 100 ${mpoint}/SCswap &> /dev/null
	start-stop-daemon --background --start \
  	--exec ${SIADEXEC} -u ${SIAUSR} --make-pidfile --pidfile ${PIDFILE} \
	-- --modules cghrtw --sia-directory ${SCDIR} --rpc-addr ${sc_rpc_port} --host-addr ${sc_host_port}
  	eend $?
}

stop() {
        ebegin "Stopping siad daemon"
	# If we want to restup SC, just remove it from startup NOW as it will fail to stop later
	# we don't source the parameters file directly
	setupsc=$(grep -P 'setup_sc\s*=\.*' /boot/storij/parameters.txt |cut -d '=' -f 2)
	if [ ${setupsc} == "1" ] || [ ${setupsc} == "2" ]; then
	     	rc-update del siad default
		rc-update del SCWalletUnlocker default
	fi
	# ensuring final stop
        start-stop-daemon --stop --exec ${SIACEXE} stop \
	--pidfile ${PIDFILE} -d ${SCDIR} || siac stop
        eend $? "Failed to stop siad daemon"
}

fstart() {
	ebegin "Starting siad daemon"
	# Sourcing parameters.txt
	/usr/local/bin/paramcopy
	# We check if ${SIAUSR} exists. If not, create it.
	if id "${SIAUSR}" >/dev/null 2>&1; then
                :
	else
	       /usr/sbin/useradd -m -G users,wheel,audio -s /bin/bash ${SIAUSR}
	fi
	if [ ! -d "${SCDIR}" ]; then
	  mkdir -p ${SCDIR}
	fi
	chown ${SIAUSR}:${SIAUSR} -R ${SCDIR}
	mpoint="$(dirname "${SCDIR}")"
	swapon -p 100 ${mpoint}/SCswap &> /dev/null
	start-stop-daemon --background --start \
  	--exec ${SIADEXEC} -u ${SIAUSR} --make-pidfile --pidfile ${PIDFILE} \
	-- --modules cghrtw --sia-directory ${SCDIR} --rpc-addr ${sc_rpc_port} --host-addr ${sc_host_port}
  	eend $?
}

fstop() {
        ebegin "Stopping siad daemon"
	# If we want to restup SC, just remove it from startup NOW as it will fail to stop later
	# we don't source the parameters file directly
	setupsc=$(grep -P 'setup_sc\s*=\.*' /boot/storij/parameters.txt |cut -d '=' -f 2)
	if [ ${setupsc} == "1" ] || [ ${setupsc} == "2" ]; then
	     	rc-update del siad default
		rc-update del SCWalletUnlocker default
	fi
	# ensuring final stop
        start-stop-daemon --stop --exec ${SIACEXE} stop \
	--pidfile ${PIDFILE} -d ${SCDIR} || siac stop
        eend $? "Failed to stop siad daemon"
}

