# app-misc/storij-tools/files/initd/ipsetter
#!/sbin/openrc-run
#
# Automatically sets a static IP from the dhcp data
#
# Copyright (c) 2017 ketetefid <ketetefid@gmail.com>
# License: GPL 3.0+
# Although storij OS has been tested by the author and is intended to be helpful,
# it comes with ABSOLUTELY NO WARRANTY, to the extent permitted by applicable law.

description="Automatically sets a static IP for storij according to the parameter int_ip"

depend() {
	need localmount copy_parameters
	after modules sysctl root expander
}

ip_setter() {
	# If net.${net_dev} is up through dhcp, we skip setting dhcp again and restarting net.${net_dev}
	#if ! grep -q config_${net_dev}[[:space:]]*=[[:space:]]*"dhcp" /etc/conf.d/net.${net_dev}; then
	#   if ! echo config_${net_dev}='"dhcp"' > /etc/conf.d/net.${net_dev} && /etc/init.d/net.${net_dev} restart; then
	#	eerror "Failed to get an IP from the router. Please check the connections."
	#	#eerror "The system will reboot in 20 seconds."
	#	#nohup bash -c 'reboot' &>/dev/null&
        #        return 1
	#   fi
	#fi

	einfo "Restarting the network device. Please wait..."
	echo config_${net_dev}='"dhcp"' > /etc/conf.d/net.${net_dev}
	/etc/init.d/net.${net_dev} restart &> /dev/null
		   
        # We wait untill we get an IP (wlan0 takes time)
	timeout=25
	i=0
	while ! [[ $(ifconfig | grep -A 1 ${net_dev} | tail -1 | awk '{print $2}') =~ ^[0-9]+\.[0-9]+\.[0-9]+\.[0-9]+$ ]]; do
           sleep 1
           let i++
           if (( $i > $timeout )); then
	      eerror "Check: failed to get an IP from the router. Please check the connections."
	      return 1
           fi
        done

	dnsaddr=`grep -i /etc/resolv.conf -e nameserver | awk '{print $2}'`
	gateway=$(/bin/ip route | awk '/default/ { print $3 }')
	ip_dhcp=`ifconfig | grep -A 1 ${net_dev} | tail -1 | awk '{print $2}'`
	netmask=`ifconfig | grep -A 1 ${net_dev} | tail -1 | awk '{print $4}'`
	if [ ${int_ip} == "auto" ]; then
		 ip_static=${ip_dhcp%.*}.133
	else
		 last_octet=`echo $int_ip | cut -d . -f 4`
		 ip_static=${ip_dhcp%.*}.$last_octet
	fi
	echo "config_${net_dev}="'"'"$ip_static netmask $netmask"'"' > /etc/conf.d/net.${net_dev}
	echo "routes_${net_dev}="'"'"default via $gateway"'"'>> /etc/conf.d/net.${net_dev}
	echo "dns_servers_${net_dev}="'"'"$dnsaddr"'"' >> /etc/conf.d/net.${net_dev}

	# sourcing parameters into rc.conf
	/usr/local/bin/paramcopy

	# wait for the chaos to settle down :))
	sleep 1
	
	einfo "Assigning the new IP..."
	/etc/init.d/net.${net_dev} restart
        # We wait untill we get an IP (wlan0 takes time)
        while ! [[ $(ifconfig | grep -A 1 ${net_dev} | tail -1 | awk '{print $2}') =~ ^[0-9]+\.[0-9]+\.[0-9]+\.[0-9]+$ ]]; do
           sleep 1
        done

	# we enable and start the services in /boot/storij/initlist again
	list=`cat /boot/storij/initlist`
	for init in $list
	do
		rc-update add $init default
		/etc/init.d/$init start
	done

	# If everything was OK, we would proceed. Otherwise, abort.
	if [ $? -ne 0 ]; then
		eerror "ipsetter failed. The process will be repeated."
		return 1
	fi		

	# We delete the initlist created by copy_parameters here
	echo "" > /boot/storij/initlist &> /dev/null

	# Substituting int_ip with the new, correct value
	sed -i 's,^\([ ]*int_ip[ ]*=\).*,\1'$ip_static',g' /boot/storij/parameters.txt
	
	# Deleting ipsetter from the startup as the IP information is now fully
	# written to /etc/conf.d/net.${net_dev} and we have started and enabled the services.
	rc-update del ipsetter default

	return 0
}

start() {
     	ebegin "Setting a static IP automatically"
   	ip_setter
	eend $?
}

