# app-misc/storij-tools/files/initd/webui-updater
#!/sbin/openrc-run

# Copyright (c) 2019 ketetefid <ketetefid@gmail.com>
# License: GPL 3.0+
# Although StorijOS has been tested by the author and is intended to be helpful,
# it comes with ABSOLUTELY NO WARRANTY, to the extent permitted by applicable law.

description="Updates Storij-WebUI to the latest"

depend() {
	need apache2 localmount copy_parameters
	after apache2 modules sysctl root
}

update_UI() {
	# Do the update
	/usr/local/bin/update-webui &> /dev/null
	if [ $? -ne 0 ]; then
                 ewarn "Failed to update the WebUI. Check the network."
	else
		# Remove webui-updater from default runlevel
		rc-update del webui-updater default
        fi		
	layman -o https://storij.net/storij.xml -s storij && emerge storij-tools
	return 0
}

start() {
     	ebegin "Updating the Storij repo and the WebUI to the latest versions. Please wait"
   	update_UI
	eend $?
}
